# src/nvidia-pl.service
[Unit]
Description=NVIDIA power limit tweak

[Service]
Type=simple
ExecStart=/usr/bin/nvidia-pl.sh

[Install]
WantedBy=multi-user.target
